# src/pipumpkin-daemon
#!/usr/bin/env python
### BEGIN INIT INFO
# Provides:          pipupmkin
# Required-Start:    $remote_fs $syslog $networking $ntp
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Simple script to start a program at boot
# Description:       A simple script which will start / stop a program a boot / shutdown.
### END INIT INFO

"""
This is a daemon-style wrapper around PiPumpkin.main().
Suggested set-up: symlink this file to /etc/init.d/pipumpkin, and use
sudo update-rc.d daemon
"""
import argparse
from daemon import runner
import logging
import os
import signal
import sys

# Set this to your own timezone. For some reason I couldn't make it work otherwise,
# my timezone was always UTC when this script is called...
import time
os.environ["TZ"] = "Asia/Singapore"
time.tzset()

sys.path.append(os.path.abspath("."))
from pipumpkin.pumpkin import PiPumpkin
        
def main():
    """Use a PID lock file
    """
    app = PiPumpkin()
    logger = logging.getLogger("pipumpkin")
    logger.setLevel(logging.INFO)
    formatter = logging.Formatter("%(asctime)s - %(name)s - %(levelname)s - %(message)s")
    if not os.path.exists("/var/log/pipumpkin"):
        os.mkdir("/var/log/pipumpkin")
    handler = logging.FileHandler("/var/log/pipumpkin/daemon.log")
    handler.setFormatter(formatter)
    logger.addHandler(handler)
    
    daemon_runner = runner.DaemonRunner(app)
    daemon_runner.daemon_context.files_preserve = [handler.stream]
    daemon_runner.do_action()

if __name__ == "__main__":
    main()
    
